# backup/crontabs/root
0 16 * * * /backup.sh >> /backups/backup.log 2>&1
